# 000-src/env/roles/docker/templates/docker-clean.timer.j2
[Unit]
Description=Scheduled run of docker-clean.service

[Timer]
OnCalendar={{ docker_clean_schedule }}
